# custom/files/etc/uci-defaults/sing-box.sh
#!/bin/sh

# Enable sing-box
uci set sing-box.main.enabled=1
